# qisCheck.timer
[Unit]
Description=Run HTWK QIS check periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

